# lab1/server-launcher.service
[Unit]
Description=Launch server at boot

[Service]
ExecStart=/etc/opt/webapp/app
Restart=always

[Install]
WantedBy=multi-user.target