# lab1/server-poller.service
[Unit]
Description=Poll server

[Service]
ExecStart=/etc/opt/webapp/poll_server.sh

[Install]
WantedBy=multi-user.target